# devops/ansible/roles/girder/templates/4_daemon/upstart.girder.conf.j2
description "Girder server"

start on startup
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 20 5

script

    cd {{ girder_path }}
    exec python -m girder > /dev/null 2>&1

end script
